# ansible/files/spacewiki-io.service
[Unit]
Name=spcewiki.io

[Service]
User=spacewiki
Group=spacewiki
WorkingDirectory=/srv/spacewiki/io-app/
ExecStart=/srv/spacewiki/virtualenv/bin/gunicorn spacewiki_io.wsgi
Requires=postgresql.service
After=postgresql.service
Before=nginx.service

[Install]
WantedBy=default.target
